# crontab
0 1,5,9,13,17,21 * * * . /a/env && /usr/local/bin/python3 -m tweetbot >> /tmp/log 2>&1
